# etc/periodic/weekly/900.certbot-renew
#!/bin/sh
#
# Define these variables in /etc/periodic.conf.local
#
# weekly_certbot_enable="YES"
#

if [ -r /etc/defaults/periodic.conf ]; then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

: ${weekly_certbot_cmd=/usr/local/bin/certbot}

case "$weekly_certbot_enable" in
[Yy][Ee][Ss])
    /usr/sbin/service nginx stop
    $weekly_certbot_cmd renew
    /usr/sbin/service nginx start
    ;;
esac
